# vscoded.service
[Unit]
Description=Run the VS code tunnel

[Service]
Type=oneshot
ExecStart=vscoded.sh

[Install]
WantedBy=multi-user.target